# example-files/chatwire-a.timer
[Unit]
Description=delayed-cw-a

[Timer]
OnBootSec=6sec

[Install]
WantedBy=basic.target
