# example-files/chatwire-a.timer
[Unit]
Description=delayed-cw-a

[Timer]
OnBootSec=6sec
Unit=chatwire-a.service

[Install]
WantedBy=basic.target
